# bin/systemd/petclinic.service
[Unit]
Description=Petclinic-Service
After=network.target syslog.target

[Service]
Type=forking

User=ec2-user
Group=ec2-user
ExecStart=/usr/local/bin/servicectl.sh petclinic startup
ExecStop=/usr/local/bin/servicectl.sh petclinic shutdown

RestartSec=10
Restart=always

[Install]
WantedBy=multi-user.target
